# robo/docs/PRODUTO/modulos/robo_matrix/referencia_robo_v2/codigo/dashboard_comercial/gerenciador_vendas/integracoes/systemd/robovendas-sync-clientes.timer
[Unit]
Description=Sincronizar Clientes Hubsoft a cada 1 minuto

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5
Persistent=true

[Install]
WantedBy=timers.target
